# systemd/dubstatsdiscover.timer
[Unit]
Description=Run the DubStats discoverer every week.

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target